# roles/acme/templates/lego-renew-notify.service.j2
[Unit]
# {{ ansible_managed }}
# Se dispara vía OnFailure= desde lego-renew.service cuando falla la renovación de un cert.
# Sin esto, un token expirado/revocado o una caída de deSEC dejaría que los certs venzan
# en silencio hasta que cada app tire errores de TLS.
Description=Alert that Let's Encrypt renewal failed

[Service]
Type=oneshot
ExecStart={{ acme_notify_script }} "Cert renewal FAILED on %H" "Certs may expire. Inspect: journalctl -u lego-renew.service -n 50"
