# roles/acme/templates/lego-renew.service.j2
[Unit]
# {{ ansible_managed }}
Description=Renew Let's Encrypt certificates via lego (DNS-01 / {{ acme_dns_provider }})
Wants=network-online.target
After=network-online.target
# Alerta (journal + ntfy opcional) si falla una renovación, para que los certs nunca venzan en silencio.
OnFailure=lego-renew-notify.service

[Service]
Type=oneshot
EnvironmentFile={{ acme_data_dir }}/desec.env
# lego v5: `run` obtiene o renueva; --renew-days lo vuelve un no-op hasta que
# queden <{{ acme_renew_days }} días, y --deploy-hook recarga solo cuando hay renovación real.
{% for cert in acme_certificates %}
ExecStart={{ acme_lego_bin }} run --accept-tos --email {{ acme_email }} --dns {{ acme_dns_provider }}{% if acme_dns_resolvers %} --dns.resolvers {{ acme_dns_resolvers }}{% endif %}{% if acme_dns_disable_rns %} --dns.propagation.disable-rns{% endif %}{% if acme_staging %} --server https://acme-staging-v02.api.letsencrypt.org/directory{% endif %} --path {{ acme_data_dir }}{% for d in cert.domains %} --domains "{{ d }}"{% endfor %} --renew-days {{ acme_renew_days }} --reuse-key --deploy-hook "/usr/bin/systemctl reload {{ acme_reload_service }}"
{% endfor %}
{% if acme_manage_selinux and (ansible_facts['selinux'].status | default('disabled')) == 'enabled' %}
# Re-etiqueta los certs recién escritos para que {{ acme_reload_service }} pueda leerlos.
ExecStartPost=/sbin/restorecon -R {{ acme_data_dir }}
{% endif %}
