# roles/acme/templates/lego-renew.timer.j2
[Unit]
# {{ ansible_managed }}
Description=Daily lego certificate renewal check

[Timer]
OnCalendar=daily
# Distribuye la carga y evita golpear a la CA de ACME en el mismo instante todos los días.
RandomizedDelaySec=3600
# Recupera la ejecución si la máquina estaba apagada en el horario programado.
Persistent=true

[Install]
WantedBy=timers.target
