# roles/cloud_torrent/templates/cloud-torrent.service.j2
[Unit]
Description=Cloud Torrent - Web-based torrent client
Documentation=https://github.com/jpillora/cloud-torrent
After=network.target

[Service]
Type=simple
User={{ cloud_torrent_user }}
Group={{ cloud_torrent_group }}
WorkingDirectory={{ cloud_torrent_working_dir }}

ExecStart={{ cloud_torrent_install_dir }}/cloud-torrent \
  --port {{ cloud_torrent_port }} \
  --host {{ cloud_torrent_host }} \
  --title "{{ cloud_torrent_title }}" \
  --config-path {{ cloud_torrent_config_path }} \
  --log{% if cloud_torrent_enable_auth %} \
  --auth {{ cloud_torrent_auth_user }}:{{ cloud_torrent_auth_password }}{% endif %}


Environment="CLOUD_TORRENT_DOWNLOADS={{ cloud_torrent_downloads_dir }}"

Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ cloud_torrent_working_dir }} {{ cloud_torrent_downloads_dir }}

[Install]
WantedBy=multi-user.target
