# roles/backup/templates/backup-maintenance.service.j2
[Unit]
# {{ ansible_managed }}
Description=restic repository maintenance (prune + check)
RequiresMountsFor={{ storage_backup_disk }}
{% if backup_notify_enabled %}
OnFailure=backup-notify@%n.service
{% endif %}

[Service]
Type=oneshot
ExecStart={{ backup_maintenance_script }}
Nice=10
IOSchedulingClass=idle
