# roles/backup/templates/backup.service.j2
[Unit]
# {{ ansible_managed }}
Description=restic backup of irreplaceable data to {{ storage_backup_disk }}
# Dependencia de montaje best-effort; el script además verifica los mounts de
# forma explícita porque RequiresMountsFor no es confiable para user units.
RequiresMountsFor={{ storage_backup_disk }} {{ storage_appdata_disk }}
{% if backup_notify_enabled %}
OnFailure=backup-notify@%n.service
{% endif %}

[Service]
Type=oneshot
ExecStart={{ backup_run_script }}
# Ser amable: los backups nunca deben dejar sin recursos a las apps.
Nice=10
IOSchedulingClass=idle
