# roles/orcaslicer/templates/orcaslicer.container.j2
[Unit]
Description=OrcaSlicer Web 3D Slicer
After=network-online.target
Wants=network-online.target

[Container]
Image={{ orcaslicer_image }}
ContainerName=orcaslicer
PublishPort={{ orcaslicer_host }}:{{ orcaslicer_port }}:3001

# Volume de config (home dir del container). Sin :Z acá — selinux.yml etiqueta
# este dir como container_file_t (una etiqueta persistente); agregar :Z aplicaría
# una categoría MCS por container que entraría en conflicto con ese etiquetado.
Volume={{ orcaslicer_config_dir }}:/config

Environment=PUID={{ orcaslicer_uid }}
Environment=PGID={{ orcaslicer_uid }}
Environment=TZ={{ orcaslicer_timezone }}
Environment=CUSTOM_USER={{ orcaslicer_custom_user }}
Environment=PASSWORD={{ orcaslicer_password }}

# Fuerza encoding por CPU (la GPU no está soportada en este hardware).
Environment=SELKIES_USE_CPU=true
Environment=DISABLE_DRI3=true

[Service]
Restart=on-failure
RestartSec=30
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target default.target
