# roles/monitoring/templates/uptime-kuma.container.j2
# Managed by Ansible - DO NOT EDIT MANUALLY
# App de un solo container: una unit .container simple (antes era un pod .kube con
# un container). El etiquetado SELinux del volumen de datos lo maneja el
# selinux.yml del role (container_file_t), así que acá no se pide relabel :z/:Z.

[Unit]
Description=Uptime-Kuma Monitoring
Documentation=https://github.com/louislam/uptime-kuma
After=network-online.target
Wants=network-online.target

[Container]
Image={{ monitoring_image }}
ContainerName=uptime-kuma

# Comparte el network namespace del host para que el monitor pueda alcanzar cada
# servicio directamente en el loopback del host (127.0.0.1:<port>): los servicios
# bindean solo a localhost, algo que un container rootless en bridge no puede
# alcanzar. Uptime-Kuma queda fijado a {{ monitoring_host }}:{{ monitoring_port }} vía las env vars de
# abajo, así que sigue escuchando solo en loopback y queda detrás de NGINX (status.<domain>).
Network=host
Environment=UPTIME_KUMA_HOST={{ monitoring_host }}
Environment=UPTIME_KUMA_PORT={{ monitoring_port }}

Volume={{ monitoring_data_dir }}:/app/data

# Los monitors de tipo "ping" (2.x) invocan /usr/bin/ping, un binario con
# cap_net_raw=ep. Podman no incluye CAP_NET_RAW en el bounding set por default, así que
# el exec del binario falla con EPERM ("spawn EPERM" en el log) aun corriendo como root.
# Se agrega la capability para que los monitors ICMP funcionen. (En 1.x andaba porque
# usaba ICMP no privilegiado; 2.x cambió a spawnear el binario ping.)
AddCapability=CAP_NET_RAW

[Service]
Restart=on-failure
RestartSec=30
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target
