# roles/service_maintenance/templates/service-watchdog.service.j2
[Unit]
Description={{ service_maintenance_unit_description }}
After=network-online.target AdGuardHome.service
Wants=network-online.target

[Service]
Type=oneshot

{% for service in service_maintenance_restart_services %}
# Only restart {{ service }} if it is NOT healthy — a blind restart on every
# boot drops LAN DNS/DHCP for no reason when the service came up fine. Healthy =
# unit active AND (dig absent OR it answers DNS on 127.0.0.1). The DNS probe
# assumes the service answers on 127.0.0.1:53 (true for AdGuard).
ExecStart=/bin/sh -c 'svc={{ service }}; if systemctl is-active --quiet "$svc" && { ! command -v dig >/dev/null 2>&1 || dig +time=2 +tries=1 @127.0.0.1 health-check.invalid >/dev/null 2>&1; }; then echo "$svc healthy, no restart"; else echo "$svc unhealthy, restarting"; systemctl restart "$svc"; fi'
{% endfor %}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ service_maintenance_unit_name }}
