# roles/service_maintenance/templates/service-watchdog.service.j2
[Unit]
Description={{ service_maintenance_unit_description }}
After=network-online.target AdGuardHome.service
Wants=network-online.target

[Service]
Type=oneshot

{% for service in service_maintenance_restart_services %}
ExecStart=/usr/bin/systemctl restart {{ service }}
{% endfor %}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ service_maintenance_unit_name }}
